# supervisord.conf
[program:logstash]
command=/opt/logstash/bin/logstash -f /etc/logstash/logstash.conf
user=logstash